# services/recorder/finish
#!/bin/bash
echo 'exited owntracks recorder'
